# tg-ai-moderator.service
[Unit]
Description=Telegram AI Moderator Bot
After=network.target

[Service]
User=tg-bot-user
Group=tg-bot-group
Type=simple
WorkingDirectory=/opt/tg-bot
EnvironmentFile=/opt/tg-bot/.env
ExecStart=/opt/tg-bot/tg-bot --db-file /var/lib/tg-bot/state.fjall
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target